# dnsmasq.ini
[program:dnsmasq]
command = /usr/sbin/dnsmasq -k
user = root
autostart = true
autoresart = true
stderr_logfile = NONE
stdout_logfile = NONE
startsecs = 3