# shadowsocksr.ini
[program:shadowsocksr]
command = python /usr/local/shadowsocksr/server.py m
user = shadowsocksr
autostart = true
autoresart = true
stderr_logfile = /var/log/shadowsocksr.log
stdout_logfile = /var/log/shadowsocksr.log
startsecs = 3
